# cmd/chihaya/deb/etc/sv/chihaya/run
#!/bin/bash

exec 2>&1
exec /usr/share/chihaya/chihaya --config /etc/chihaya/config.yaml --debug
